# groovyarcade/airootfs/etc/profile.d/aliases.sh
#!/bin/bash

alias ls='ls --color=auto'
alias startx="startx -- vt7"
alias ll="ls -l"
