# groovyarcade/airootfs/home/arcade/.config/systemd/user/udiskie.service
[Unit]
Description=udiskie mounts drives when plugged in

[Service]
ExecStart=/usr/bin/udiskie --no-notify

[Install]
WantedBy=default.target
